# backend/src/ops/systemd/ag.service
# ag.service - Agentic RAG server (user service)
# Recommended for workstations: install as a user service under
# ~/.config/systemd/user/ag.service
# This unit uses per-user paths and does not set User/Group.

[Unit]
Description=Agentic RAG API Server (user)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Working directory for relative paths like "documents" (per-user XDG)
WorkingDirectory=%h/.local/share/ag

# Environment: per-user env file (XDG config)
EnvironmentFile=-%h/.config/ag/ag.env

# Example inline overrides (optional)
# Environment=RUST_LOG=info
# Environment=BACKEND_HOST=127.0.0.1
# Environment=BACKEND_PORT=3010
# Environment=SKIP_INITIAL_INDEXING=true
# Environment=TRUST_PROXY=false
# Environment=RATE_LIMIT_ROUTES_FILE=%h/.config/ag/rl-routes.yaml

# Binary path (per-user install)
ExecStartPre=/bin/bash -c 'for i in $(seq 1 10); do ss -tlnp | grep -q :3010 || exit 0; sleep 1; done; fuser -k 3010/tcp 2>/dev/null; sleep 1'
ExecStart=%h/.local/bin/ag
StandardOutput=journal
StandardError=journal

Restart=on-failure
RestartSec=2
KillMode=control-group
TimeoutStopSec=15

# Resource/FD limits
LimitNOFILE=65536

# Hardening (optional; relax as needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=default.target
