# systemd/ag.service.d/falkordb.conf
# Drop-in: ag.service depends on falkordb.service.
# Installed by install-linux.sh unless --no-falkordb was passed.
[Unit]
Wants=falkordb.service
After=falkordb.service
