# systemd/ag.service.d/stack.conf
# Drop-in: ag.service depends on ag-stack.service (docker compose stack).
# Installed by install-linux.sh unless --no-stack was passed.
[Unit]
Wants=ag-stack.service
After=ag-stack.service
